# deploy/systemd/passwall-sub-panel.service
[Unit]
Description=Passwall-Sub-Panel
Documentation=https://github.com/KazuhaHub/Passwall-Sub-Panel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=psp
Group=psp
WorkingDirectory=/opt/psp
# Place secrets and DSN here (chmod 600):
#   PSP_SECRET_KEY=...
#   PSP_MYSQL_DSN=user:pass@tcp(127.0.0.1:3306)/psp?...
#   PSP_JWT_SECRET=...
EnvironmentFile=/opt/psp/.env
ExecStart=/opt/psp/psp --config /opt/psp/config.yaml
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/psp/data
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
